# sahra-ai.service
[Unit]
Description=Sahra AI — Quran-based commerce assistant
After=network.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/quranai
EnvironmentFile=-/opt/quranai/.env
ExecStart=/bin/bash /opt/quranai/start.sh
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
