# sahra-ai.service
[Unit]
Description=Sahra AI — Quran-based commerce assistant
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/sahra
EnvironmentFile=/opt/sahra/.env
ExecStart=/bin/bash /opt/sahra/start.sh
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
